# noir-gateway.service
[Unit]
Description=Noir:GateWay — Badge Elevation System
After=network.target

[Service]
Type=exec
User=root
WorkingDirectory=/root/Builds/Noir:GateWay
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=NEXT_TELEMETRY_DISABLED=1

[Install]
WantedBy=multi-user.target
